# deploy/economist-rss-refresh.timer
[Unit]
Description=Refresh private Economist full-text RSS cache every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
RandomizedDelaySec=1min
Persistent=true
Unit=economist-rss-refresh.service

[Install]
WantedBy=timers.target
