# deploy/economist-rss-refresh.timer
[Unit]
Description=Refresh private Economist full-text RSS cache every hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=15min
Persistent=true
Unit=economist-rss-refresh.service

[Install]
WantedBy=timers.target
