# deploy/economist-rss-refresh.timer
[Unit]
Description=Refresh private Economist full-text RSS cache every 10 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min
RandomizedDelaySec=2min
Persistent=true
Unit=economist-rss-refresh.service

[Install]
WantedBy=timers.target
